# supervisord/php-fpm.conf
[program:php-fpm]
command=/usr/local/php/sbin/php-fpm -F --fpm-config /usr/local/php/etc/php-fpm.conf
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0